# workshops/lab-containerizing-spring/workshop/profile.d/use-exercises-directory.sh
# Temporary hack until Educates fixed to calculate these at end of setup.

EXERCISES_DIR=exercises

TERMINAL_HOME=$HOME/$EXERCISES_DIR
export TERMINAL_HOME

EDITOR_HOME=$HOME/$EXERCISES_DIR
export EDITOR_HOME
